# examples/monitor/systemd/whoopskill-monitor.service
[Unit]
Description=whoopskill token refresh + lightweight fetch

[Service]
Type=oneshot
# Ensure your WHOOP_* env vars are set via ~/.config/environment.d/*.conf
# or via an EnvironmentFile= line below.
# EnvironmentFile=%h/.config/whoopskill.env
ExecStart=%h/.local/bin/whoop-refresh-monitor.sh
